# pokemon/deploy_dev/supervisord.dir/bk/02_game_server.ini
[program:02_game_server]  
directory=../release/
command=/usr/bin/python game_server.py game.dev.2
autostart=false  
autorestart=false
priority=103
redirect_stderr=true
stdout_logfile_maxbytes=100MB
stdout_logfile_backups=999999
