# pokemon/deploy_dev/supervisord.dir/bk/02_pvp_server.ini
[program:02_pvp_server]  
directory=../release/pvp
command=%(here)s/../release/bin/host_server -config=pvp.dev.2
autostart=false  
autorestart=false
priority=102
redirect_stderr=true

