# pokemon/deploy_dev/supervisord.dir/bk/02_storage_server.ini
[program:02_storage_server]
directory=../release/storage
command=%(here)s/../release/bin/host_server -config=storage.dev.2
autostart=false  
autorestart=false
priority=101
redirect_stderr=true

