# pokemon/deploy_dev/supervisord.dir/bk/03_pvp_server.ini
[program:03_pvp_server]  
directory=../release/pvp
command=%(here)s/../release/bin/host_server -config=pvp.dev.3
autostart=false  
autorestart=false
priority=102
redirect_stderr=true

