# pokemon/deploy_dev/supervisord.dir/bk/anti_cheat_server.ini
[program:anti_cheat_server]  
directory=../release/anti_cheat
command=%(here)s/../release/bin/anti_cheat_server -config=agent.dev.1 -limit=1
autostart=false  
autorestart=false
priority=102
redirect_stderr=true

