# pokemon/deploy_dev/supervisord.dir/crash_platform_server.ini
[program:crash_platform_server]  
directory=/mnt/pokemon/crash_platform/
command=/usr/bin/python server.py --name=crash_platform_qq
autostart=false  
autorestart=false
priority=105
redirect_stderr=true
