# pokemon/deploy_dev/supervisord.dir/disable_word_check_server.ini
[program:disable_word_check_server]
directory=/mnt/pokemon/release/disable_word_check
command=/mnt/pokemon/release/bin/disable_word_check_server 192.168.1.233:38080
autostart=false
autorestart=true
priority=101
redirect_stderr=true
