# pokemon/deploy_dev/supervisord.dir/gm_server.ini
[program:gm_server]  
directory=../release/  
command=/usr/bin/python gm_server.py gm
autostart=false  
autorestart=false
priority=105
redirect_stderr=true

