# pokemon/deploy_dev/supervisord.dir/login_server.ini
[program:login_server]
directory=%(here)s/../release/login
command=%(here)s/../release/bin/login_server -config=login.dev.1
autostart=false
autorestart=false
priority=110
redirect_stderr=true
stdout_logfile_maxbytes=100MB
stdout_logfile_backups=999999
