# pokemon/deploy_dev/supervisord.dir/mongodb.ini
[program:game_mongodb]
directory=%(here)s
command=mongod -f ./mongodb.conf --dbpath=./game_db
autostart=false
autorestart=false
startsecs=5
priority=100
redirect_stderr=true
stopwaitsecs=99999999
