# pokemon/deploy_dev/supervisord.dir/nsqadmin.ini
[program:nsqadmin]  
directory=./nsq
command=/mnt/pokemon/deploy_dev/nsq/nsqadmin --lookupd-http-address=127.0.0.1:4161
autostart=false
autorestart=true
priority=100
redirect_stderr=true

