# pokemon/deploy_dev/supervisord.dir/nsqd.ini
[program:nsqd]  
directory=./nsq
command=/mnt/pokemon/deploy_dev/nsq/nsqd --lookupd-tcp-address=127.0.0.1:4160 --snappy=true --max-msg-size=67108864
autostart=false
autorestart=true
priority=100
redirect_stderr=true

