# pokemon/deploy_dev/supervisord.dir/payment_server.ini
[program:payment_server]  
directory=../release/
command=/usr/bin/python payment_server.py payment.dev.1
autostart=false  
autorestart=false
priority=103
redirect_stderr=true
stdout_logfile_maxbytes=100MB
stdout_logfile_backups=999999
